# wine32_build/ukylin-wine/DEBIAN/prerm
#!bin/sh

umount /run/wine
systemctl disable run-wine.mount
